# templates/deluged.service.erb
[Unit]
Description=Deluge Bittorrent Client Daemon
After=network-online.target

[Service]
Type=simple
User=deluge
Group=deluge
UMask=007

ExecStart=/usr/bin/deluged -d <%= node['deluge']['daemon']['logs']['enabled'] %> <%= node['deluge']['daemon']['logs']['file'] %> <%= node['deluge']['daemon']['logs']['level'] %>

Restart=on-failure

# Configures the time to wait before service is stopped forcefully.
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
